# ed_quant_engine/quant_bot.service
[Unit]
Description=ED Capital Quant Engine Daemon
After=network.target

[Service]
Type=simple
User=root
# Change WorkingDirectory to the actual absolute path of ed_quant_engine
WorkingDirectory=/path/to/ed_quant_engine
# Use the virtual environment's Python executable
ExecStart=/path/to/ed_quant_engine/venv/bin/python /path/to/ed_quant_engine/main.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=quant_engine
Environment=TZ=Europe/Istanbul
EnvironmentFile=/path/to/ed_quant_engine/.env

[Install]
WantedBy=multi-user.target
