# ed_quant_engine/quant_bot.service
[Unit]
Description=ED Capital Quant Engine Service
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
# Adjust this path to the absolute path of your project directory
WorkingDirectory=/path/to/ed_quant_engine

# Start the docker-compose setup
ExecStart=/usr/local/bin/docker-compose up
ExecStop=/usr/local/bin/docker-compose down

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
