# ed_quant_engine/quant_bot.service
[Unit]
Description=ED Capital Quant Engine Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ed_quant_engine
ExecStart=/opt/ed_quant_engine/manage_bot.sh start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
