# ed_quant_engine/quant_bot.service
[Unit]
Description=ED Capital Quant Engine Daemon
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ed_quant_engine
ExecStart=/usr/local/bin/docker-compose up
ExecStop=/usr/local/bin/docker-compose down
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
