# ed_quant_engine/quant_bot.service
[Unit]
Description=ED Capital Quant Engine
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ed_quant_engine
ExecStart=/opt/ed_quant_engine/manage_bot.sh start
ExecStop=/opt/ed_quant_engine/manage_bot.sh stop
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
