# ed_quant_engine/quant_bot.service
[Unit]
Description=JulesEmtia Quant Engine
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/JulesEmtia/ed_quant_engine
ExecStart=/opt/JulesEmtia/ed_quant_engine/manage_bot.sh start
ExecStop=/opt/JulesEmtia/ed_quant_engine/manage_bot.sh stop
ExecReload=/opt/JulesEmtia/ed_quant_engine/manage_bot.sh restart
Restart=on-failure
RestartSec=10
Environment="TZ=Europe/Istanbul"

[Install]
WantedBy=multi-user.target
